# src/main/python/resources/alppaca
#!/bin/bash
#
# /etc/rc.d/init.d/alppaca
#
### BEGIN INIT INFO
# Provides: alppaca
# Required-Start: $network, $syslog
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Alppaca
# Description: A(mazing)LocalPrefetchProxyforAmazonCredentiAls
### END INIT INFO

alppacad $@
